# scripts/sentient-learning.service
[Unit]
Description=SentientOS Continuous Learning Service
After=network.target

[Service]
Type=simple
User=sentient
WorkingDirectory=/opt/sentientos
Environment="PATH=/opt/sentientos/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/sentientos/.venv/bin/python /opt/sentientos/rl_agent/orchestrator.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/sentientos/logs/continuous_learning.log
StandardError=append:/opt/sentientos/logs/continuous_learning.log

[Install]
WantedBy=multi-user.target